# packaging/linux/prerm
#!/bin/sh
set -eu

case "${1:-}" in
  remove|deconfigure)
    if [ -d /run/systemd/system ]; then
      systemctl disable --now ket-tunnel.service || true
    fi
    ;;
esac

exit 0
